# apartment-bot.service
[Unit]
Description=Apartment Telegram Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/apartment_bot
EnvironmentFile=/home/ubuntu/apartment_bot/.env
ExecStart=/home/ubuntu/apartment_bot/venv/bin/python main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
